# roles/kubernetes/confd/templates/confd.service.j2
[Unit]
After=etcd2.service
Description=Confd
Requires=etcd2.service

[Service]
ExecStart=/opt/bin/confd
Restart=always

[Install]
WantedBy=multi-user.target
